# robot/config/etc/systemd/system/sciveo-robot.timer
[Unit]
Description=Sciveo Robot

[Timer]
OnBootSec=0
Unit=sciveo-robot.service

[Install]
WantedBy=timers.target